# templates/launchd.media-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.media-digest</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>-m</string>
    <string>media_digest</string>
    <string>--config</string>
    <string>/path/to/Github-trend/config.toml</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/path/to/Github-trend</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/path/to/Github-trend/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/path/to/Github-trend/logs/launchd.err.log</string>
</dict>
</plist>
